# ansible/roles/rknn/templates/rkllama.service.j2
[Unit]
Description=RKLLAMA LLM Server for RK3588 NPU
After=network.target
Documentation=https://github.com/jfreed-dev/rkllama

[Service]
Type=simple
User=root
WorkingDirectory=/opt/rkllama
Environment="PATH={{ rknn_venv_path }}/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart={{ rknn_venv_path }}/bin/python3 /opt/rkllama/server.py --target_platform {{ rknn_target_platform }} --port {{ rknn_service_port }}
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
